# packages/aconnectd_0.1.0_amd64/DEBIAN/postinst
#!/bin/bash

chmod +x /opt/aconnectd/bin/aconnectd
ln -sf /opt/aconnectd/bin/aconnectd /usr/local/bin/aconnectd
systemctl enable aconnectd
systemctl start aconnectd