# packages/aconnectd_amd64/DEBIAN/prerm
#!/bin/bash

systemctl stop aconnectd
systemctl disable aconnectd